# crontab
# Run Every 04:30am
30 4 * * * /bin/bash -l -c 'backup perform -t default'